# services/danbyte-utilization.service
[Unit]
Description=Danbyte prefix-utilization alerts
After=network.target
Wants=danbyte-infra.service

[Service]
Type=oneshot
WorkingDirectory=%h/danbyte
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/danbyte/.venv/bin/python manage.py check_utilization

# Triggered every 15 minutes by danbyte-utilization.timer.
# Logs: journalctl --user -fu danbyte-utilization
